# rootfs/etc/services.d/shell-operator/run
#!/usr/bin/with-contenv bash

/shell-operator